# net-misc/openntpd/files/openntpd-3.9_p1-r3/ntpd
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

name="OpenNTPD"
command="/usr/sbin/ntpd"
command_args="${NTPD_OPTS}"

depend() {
	need net
	after ntp-client
	use dns logger
}

start_pre() {
	if [ ! -f /etc/ntpd.conf ] ; then
		eerror "Could not find /etc/ntpd.conf -- exiting."
		return 1
	fi

	if [ -x /usr/bin/getent ] ; then
		if [ "`getent passwd ntp | cut -d: -f 6`" != "${NTPD_HOME}" ] ; then
			eerror "Home directory of ntp needs to be ${NTPD_HOME}"
			eerror "Please run 'usermod -d ${NTPD_HOME} ntp'"
			return 1
		fi
	fi

	return 0
}

