# sys-apps/usbutils/files/usbutils.cron
#!/bin/sh
exec /usr/sbin/update-usbids -q
