# sys-apps/kmod/files/kmod-static-nodes
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-apps/kmod/files/kmod-static-nodes,v 1.1 2013/08/09 19:08:18 ssuominen Exp $

description="Create list of required static device nodes for the current kernel"

depend() {
	after dev-mount
	before tmpfiles.dev dev
}

start() {
	ebegin "Creating list of required static device nodes for the current kernel"
	checkpath -q -d /run/tmpfiles.d
	kmod static-nodes --format=tmpfiles --output=/run/tmpfiles.d/kmod.conf
	eend $?
}
