# media-sound/lastfmsubmitd/files/lastmp
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need mpd lastfmsubmitd
	after mpd
}

start() {
	ebegin "Starting lastmp"
	start-stop-daemon --start --quiet --chuid lastfm \
		--exec /usr/bin/lastmp
	eend $?
}

stop() {
	ebegin "Stopping lastmp"
	start-stop-daemon --stop --quiet --user lastfm \
		--exec /usr/bin/lastmp
	eend $?
}

