# media-sound/lastfmsubmitd/files/lastfmsubmitd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need localmount
}

start() {
	ebegin "Starting lastfmsubmitd"
	start-stop-daemon --start --quiet --chuid lastfm \
		--exec /usr/bin/lastfmsubmitd
	eend $?
}

stop() {
	ebegin "Stopping lastfmsubmitd"
	start-stop-daemon --stop --quiet --user lastfm \
		--exec /usr/bin/lastfmsubmitd
	eend $?
}

