# sys-fs/udev/files/udev-postmount-initd-111-r2
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-fs/udev/files/udev-postmount-initd-111-r2,v 1.1 2007/05/25 13:51:07 zzam Exp $

depend() {
	need localmount
}

start() {
	# check if this system uses udev
	[ -d /dev/.udev/ ] || return 0

	# store persistent-rules that got created while booting
	# when / was still read-only
	/lib/udev/move_tmp_persistent_rules.sh
}

stop() {
	:
}

# vim:ts=4
