# media-sound/oss/files/init.d/oss
#!/sbin/runscript

depend() {
	need bootmisc localmount
	after modules isapnp coldplug hotplug
}

start() {
	ebegin "Starting OSS"
	/usr/sbin/soundon
	/usr/sbin/savemixer -L
	eend $? "Failed to start OSS"
}

stop() {
	ebegin "Stopping OSS"
	/usr/sbin/savemixer
	/usr/sbin/soundoff
	eend $? "Failed to stop OSS"
}
