# net-misc/openvswitch/files/ovs-brcompatd
#!/sbin/runscript

depend() {
	need ovs-vswitchd
}

start() {
	ebegin "Starting Open vSwitch Linux Bridge daemon"
	start-stop-daemon --start --quiet --exec /usr/sbin/ovs-brcompatd -- --pidfile --detach ${DB_SOCKET}
	eend ${?}
}

stop() {
	ebegin "Stopping Open vSwitch Linux Bridge daemon"
	start-stop-daemon --stop --quiet --name ovs-brcompatd
	eend ${?}
}
