# net-misc/openvswitch/files/ovs-vswitchd
#!/sbin/runscript

depend() {
	need ovsdb-server
}

start() {
	ebegin "Starting Open vSwitch daemon"
	start-stop-daemon --start --quiet --exec /usr/sbin/ovs-vswitchd -- --pidfile --detach ${DB_SOCKET}
	eend ${?}
}

stop() {
	ebegin "Stopping Open vSwitch daemon"
	start-stop-daemon --stop --quiet --name ovs-vswitchd
	eend ${?}
}
