# media-sound/scrobby/files/scrobby.init
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net
	use mpd
}

start () {
	ebegin "Starting scrobby"
	start-stop-daemon --start --exec /usr/bin/scrobby
	eend $?
}

stop () {
	ebegin "Stopping scrobby"
	start-stop-daemon --stop --exec /usr/bin/scrobby || eerror "scrobby is not running"
	eend $?
}
