# sys-fs/mpd-updater/files/mpd-updater
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

start() {
	ebegin "Starting mpd-updater"
	start-stop-daemon --start --quiet --pidfile /var/run/mpd-updater.pid \
		--background --exec /usr/bin/mpd-updater
	eend $?
}

stop() {
	ebegin "Shutting down mpd-updater"
	start-stop-daemon --stop --quiet --pidfile /var/run/mpd-updater.pid
	eend $?
}
