# sys-apps/pciutils/files/pciutils.cron
#!/bin/sh
exec /usr/sbin/update-pciids -q
