# sys-apps/man-db/files/man-db.cron
#!/bin/sh
install -o man -g root -m2755 -d /var/cache/man
[ ! -x /usr/bin/mandb ] && echo "/usr/bin/mandb" does not exist && exit 1
/usr/bin/mandb -q
if [ $? -ne 0 ]
then
	echo "Error on first man attempt; attempting full rebuild."
	rm -rf /var/cache/man
	install -o man -g root -m2755 -d /var/cache/man
	exec /usr/bin/mandb -c -q
fi
