# www-servers/cherokee/files/cherokee-initd-0.99.23
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

opt="reload"

depend() {
	need net
}

start() {
	ebegin "Starting Cherokee Web Server"
	start-stop-daemon --start -b \
		--exec /usr/sbin/cherokee
	eend $?
}

stop() {
	ebegin "Stopping Cherokee Web Server"
	start-stop-daemon --stop \
		--exec /usr/sbin/cherokee
	eend $?
}

reload() {
	ebegin "Restarting Cherokee Web Server and closing all open connections"
	start-stop-daemon --stop \
		--signal USR1
	eend $?
}

