# sys-power/cpufrequtils/files/cpufrequtils-init.d-007
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-power/cpufrequtils/files/cpufrequtils-init.d-007,v 1.2 2013/05/06 11:09:09 lxnay Exp $

change() {
	ebegin "Running cpufreq-set ${*}"
	/usr/libexec/cpufrequtils-change.sh "${@}" > /dev/null
	eend $?
}

start() {
	change ${START_OPTS} -- ${SYSFS_EXTRA}
}

stop() {
	change "${STOP_OPTS}"
}
