# dev-lang/php/files/5.3.0_beta1/init.d/phpd
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

export "PHP_FCGI_CHILDREN=${PHP_FCGI_CHILDREN} PHP_FCGI_MAX_REQUESTS=${PHP_FCGI_MAX_REQUESTS} USERID=${USERID}"

checkconfig() {
    if [ ! -f ${CONFIGDIR}/php.ini ] ; then
        eerror "No ${CONFIGDIR}/php.ini file exists!"
        return 1
    fi
}

start() {
    ebegin "Checking phpd configuration"
        checkconfig || return 1
    eend $?

    ebegin "Starting phpd"
        start-stop-daemon --start --quiet --background \
            --pidfile=${PIDFILE} --make-pidfile \
            --exec=${DAEMON} --user=${USERID} --group=${GROUP} -- ${DAEMON_ARGS}
    eend $?
}

stop () {
    ebegin "Stopping phpd"
        start-stop-daemon --stop --quiet --pidfile=${PIDFILE}
    eend $?
}
