# net-misc/pooler-cpuminer/files/minerd.initd
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
        need net
}

start() {
        ebegin "Starting Miner daemon"
	start-stop-daemon \
	--start --user "${MINERD_USER}" --name minerd \
	--pidfile /var/run/minerd.pid --make-pidfile \
	--env HOME="${VARDIR}" --exec /usr/bin/minerd \
	--nicelevel "${NICELEVEL}" \
	--background \
	--wait 2000 \
	-- ${MINERD_OPTS}
	eend $?

}

stop() {
        ebegin "Stopping Miner daemon"
	start-stop-daemon --stop --user "${MINERD_USER}" \
	--name minerd --pidfile /var/run/minerd.pid \
	--wait 30000 \
	--progress
	eend $?
}

