# net-misc/dhcpcd/files/dhcpcd.initd-r2
#!/sbin/runscript
# Copyright 2008-2010 Piotr Karbowski <piotr.karbowski@gmail.com>
# All rights reserved. Released under the 2-clause BSD license.

pidfile=/var/run/dhcpcd.pid
name="DHCP Client Daemon"

depend()
{
	provide dhcp-client
	need localmount
	use logger
	before dns
}

start() {
	ebegin "Starting ${name}"
	start-stop-daemon --start --quiet --pidfile ${pidfile} \
		--exec /sbin/dhcpcd -- -q
}
